# root/etc/services.d/satip/run
#!/usr/bin/with-contenv bash

cd /app/satip || exit

IFS=" " read -r -a RUN_ARRAY <<< "$RUN_OPTS"

exec \
        ./minisatip "${RUN_ARRAY[@]}" -f -x 8875
