# root/etc/services.d/satip/run
#!/usr/bin/with-contenv bash

cd "${SATIP_ROOT}"

./minisatip -x 8875
