# root/etc/services.d/satip/run
#!/usr/bin/with-contenv bash

cd /app/satip

exec \
	./minisatip "${RUNCMD}" -f -x 8875
