# root/etc/services.d/satip/run
#!/usr/bin/with-contenv sh
cd /app/satip

./minisatip -x 8875
